# crontab
#!/bin/bash

# Wait a random amount of time between 0–30 minutes (0–1800 seconds)
sleep $((RANDOM % 1800))

# Run the actual DNS noise script
/your-directory/dns_noise.sh
